# systemd-user/blueman-applet.service
[Install]
WantedBy=desktop-session.target

[Service]
ExecStart=/usr/bin/blueman-applet

[Unit]
Description=Blueman applet
Requires=tray.target
