# systemd-user/tray.target
[Unit]
Description=System Tray
WantedBy=desktop-session.target
