# systemd-user/xsettingsd.service
[Install]
WantedBy=desktop-session.target

[Service]
ExecStart=/usr/bin/xsettingsd
Restart=on-abort

[Unit]
Description=xsettingsd
